# ide/org.eclipse.scada.configuration.world.lib/src/org/eclipse/scada/configuration/world/lib/deployment/startup/templates/systemd/driver.service
[Unit]
Description=Eclipse SCADA @@driverName@@ driver
After=network.target

[Service]
ExecStart=/usr/bin/scada.driver.@@driverName@@.launcher
StandardOutput=null
User=neoscada

[Install]
WantedBy=multi-user.target
