# ide/org.eclipse.scada.configuration.world.lib/src/org/eclipse/scada/configuration/world/lib/deployment/startup/templates/systemd/p2.service
[Unit]
Description=Equinox application @@appName@@
After=network.target

[Service]
ExecStart=/usr/bin/scada.app.@@appName@@.launcher
StandardOutput=null
User=neoscada

[Install]
WantedBy=multi-user.target
