# ide/org.eclipse.scada.configuration.world.lib/src/org/eclipse/scada/configuration/world/lib/deployment/startup/templates/upstart/app.upstart.conf
# start $$appName$$

author "$$authorName$$"
description "Equinox application $$appName$$"

start on (filesystem and runlevel [2345])
stop on runlevel [016]

respawn

pre-start script
	[ -d ~neoscada/"$$appName$$" ] || stop && exit 0
end script

script
    . /etc/default/neoscada
    exec su neoscada -s /bin/bash -c "cd ~neoscada/\"$$appName$$\"; exec screen -D -m -S \"$$appName$$\" ./launcher"
end script
