# ide/org.eclipse.scada.configuration.world.lib/src/org/eclipse/scada/configuration/world/lib/deployment/startup/templates/upstart/app.upstart.sf.conf
# start $$appName$$

author "$$authorName$$"
description "Equinox application $$appName$$"

start on (filesystem and runlevel [2345])
stop on runlevel [016]

respawn

pre-start script
    # Work around a race condition in screen creating "/var/run/screen"
	[ -d /var/run -a ! -d /var/run/screen ] && install -m 0775 -o root -g utmp -d /var/run/screen 
	[ -d ~neoscada/"$$appName$$" ] || stop && exit 0
end script

script
    . /etc/default/neoscada
    exec su neoscada -s /bin/bash -c "cd ~neoscada/\"$$appName$$\"; exec screen -D -m -S \"$$appName$$\" ./launcher"
end script
