# ide/org.eclipse.scada.configuration.world.lib/src/org/eclipse/scada/configuration/world/lib/deployment/startup/templates/upstart/driver.upstart.conf
author "$$authorName$$"
description "Eclipse NeoSCADA $$driverName$$ driver"

env driverName="$$driverName$$"

start on (filesystem and runlevel [2345])
stop on runlevel [016]

respawn

pre-start script
	. /etc/default/neoscada
	. esFunctions
	[ -f "$NEOSCADA_CONFIGURATION_ROOT/drivers/$driverName/exporter.xml" ] || stop && exit 0
end script

script
    . /etc/default/neoscada
    exec esDriverExporter "$driverName" -attach
end script